# scripts/chef/s3-backup/files/default/cron.d/backup-mysql
0 9 * * 0 root /usr/local/bin/backup-mysql.sh 2>&1 | logger -p local3.info
0 0 * * 1-5 root /usr/local/bin/backup-mysql.sh 2>&1 | logger -p local3.info
