# scripts/chef/s3-backup/files/default/cron.d/backup-retention
0 0 * * * backup /usr/local/bin/backup-retention.sh

